# crontab
*/1 * * * * /app/testFinal.py
*/1 * * * * touch /app/cron_exist.cron
