# crontab
* * * * * python testFinal.py >> /app/cron.log 2>&1
00 08 * * * python generateurModele.py >> /app/cron.log 2>&1

