# crontab
*/1 * * * * python3 -u /app/generateurModele.py
*/1 * * * * touch /app/cron_exist.cron
