# crontab
*/1 * * * * python3 -u /app/generateurModele.py

