# crontab
00 08 * * * python generateurModele.py
* * * * * python testFinal.py

