# crontab
*/1 * * * * python3 -u /app/testFinal.py
*/1 * * * * touch /app/cron_exist.cron
