# crontab
*/1 * * * * python3 /app/generateurModele.py

