# crontab
*/1 * * * * python3 -u /app/testFinal.py &> /app/cron_py.log
*/1 * * * * touch /app/cron_exist.cron
