# crontab
* * * * * sudo python3 /app/testFinal.py >> /app/cron.log 2>&1
00 08 * * * sudo python3 /app/generateurModele.py >> /app/cron.log 2>&1

